# proxy/scripts/letsencrypt.cron
27 4 * * * /bin/bash /opt/scripts/letsencrypt-renew.sh
